# docker/php-flex/files/etc/supervisord.conf
[supervisord]
nodaemon=true
umask=022
pidfile=/var/run/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0
user=root

[program:messenger-consume]
command=/srv/bin/console messenger:consume async --time-limit=1800 -vvv
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
exitcodes=0
startsecs=1
startretries=24
user=www-data
stopsignal=INT
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
